# setup/ttyd.service
[Unit]
Description=ttyd - Terminal in browser
After=network.target

[Service]
Type=simple
User=pds
WorkingDirectory=/home/pds
Environment="HOME=/home/pds"
Environment="TERM=xterm-256color"
# Persistent tmux session - all browser tabs share the same session
# Use tmux windows for multiple terminals
# Theme: Midnight Command (#060a14 background, #ff8c00 cursor)
# Uses wrapper script to avoid systemd JSON escaping issues
ExecStart=/home/pds/boomshakalaka/scripts/start_ttyd.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
